# X220/source/initramfs-update.service
[Unit]
Description=Copy initramfs Kernel to EFI System Partition

[Service]
Type=oneshot
ExecStart=/usr/share/megavolts/copy_efistub.sh
